# .config/crontab/ariel.tom.crontab
SHELL=/bin/bash
LC_ALL=en_US.UTF-8
PATH=/usr/local/bin:/usr/bin:/bin:/usr/local/sbin:/usr/sbin:/sbin:/Users/tom/bin

# min  hr  dom  mon  dow     command
*/5   *   *    *    *       touch /Volumes/RAID10/.nosleep > /dev/null 2>&1
#*/5   *   *    *    *       diskutil mount RAID10 > /dev/null 2>&1
#  5    0   *    *    0,2,4   ~/bin/sync-annex.sh content && ~/bin/sync-annex.sh fsck
#  5    0   *    *    1,3,5   ~/bin/sync-annex.sh content && ~/bin/sync-annex.sh fsck remotes
#  5    0   *    *    6       ~/bin/sync-annex.sh all

#@daily         ~/bin/rsyncBackup_daily
@daily          ~/bin/getMovieList
@daily          ~/bin/updateHostsFile.sh runfromcron
#@daily          for dir in $HOME/annex/*/ ; do cd "${dir}" && git annex watch ; done
#@daily          for dir in /Volumes/RAID10/annex/*/ ; do cd "${dir}" && git annex watch ; done
#@daily          for dir in $HOME/annex/*/ ; do cd "${dir}" && git annex list > "${dir%/}.txt" ; done
@daily          for dir in $HOME/Documents/FileThis\ Fetch/*/ ; do rmdir "${dir}" ; done

#@weekly        ~/bin/rsyncBackup_weekly
@weekly         ~/bin/homebrewBackup_weekly.sh
@weekly         rsync -qcltz --block-size=8192 --partial --progress --exclude=emeritus-* --exclude=removed-* keyring.debian.org::keyrings/keyrings/* ~/.debian-keyring

#
00	23	*	*	*	/Users/tom/Library/Application\ Support/SuperDuper\!/Scheduled\ Copies/SDUtil -i || open file:///Users/tom/Library/Application\%20Support/SuperDuper\%21/Scheduled\%20Copies/Smart\%20Update\%20c\%20from\%20a.sdsp/Copy\%20Job.app
#
